# nas-backup.timer
[Unit]
Description=Run NAS backup nightly

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
